# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:node-server]
command=node /app/server/server.js
directory=/app
environment=NODE_ENV="production",PORT="3001",PYTHONIOENCODING="utf-8"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=2
startretries=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=2
startretries=5

[program:fund-sidecar]
command=python -m uvicorn engine_fund.api.sidecar:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app
environment=PYTHONPATH="/app",PYTHONIOENCODING="utf-8",REDIS_URL="%(ENV_REDIS_URL)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=3
startretries=5

[program:celery-worker]
command=celery -A engine_fund.api.celery_app:celery_app worker --loglevel=info --pool=solo --concurrency=1 --max-tasks-per-child=10
directory=/app
environment=PYTHONPATH="/app",PYTHONIOENCODING="utf-8",REDIS_URL="%(ENV_REDIS_URL)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
startretries=5
; Allow Celery up to 3 minutes for warm shutdown so an in-flight simulation
; can either finish or have its message re-queued cleanly via task_acks_late.
stopwaitsecs=180
; SIGTERM = warm shutdown (Celery finishes current task before exiting).
stopsignal=TERM
; Kill the entire process group so any child Python processes also exit
; (and don't leave zombie matplotlib/numpy workers behind).
stopasgroup=true
killasgroup=true
